# docker/config/websockets/websockets.ini
[program:websockets]
command=/usr/local/bin/php /var/www/html/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
